# rattler433/rootfs/etc/services.d/rtl_433_devices/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Rattler 433
# Creates MQTT devices in Home Assistant
# ==============================================================================

/usr/bin/python3 /rattler/device_creator/main.py
