# rattler433/rootfs/etc/services.d/rtl_433_devices/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Rattler 433
# Creates MQTT devices in Home Assistant
# ==============================================================================
# Declarations for MQTT
declare ha_discovery_topic
declare host
declare password
declare port
declare username

# Declarations for devices
declare channel
declare id
declare manufacturer
declare model
declare name
declare type
declare uid

if ! bashio::services.available "mqtt"; then
  bashio::log.info "No internal MQTT service found"
  if bashio::config.has_value "mqtt_host"; then
    host=$(bashio::config 'mqtt_host')
  else
    bashio::log.info "No external MQTT service defined"
  fi
  if bashio::config.has_value "mqtt_port"; then
    port=$(bashio::config 'mqtt_port')
  else
    bashio::log.info "No external MQTT service defined"
  fi
  if bashio::config.has_value "mqtt_user"; then
    username=$(bashio::config 'mqtt_user')
  else
    bashio::log.info "No external MQTT service defined"
  fi
  if bashio::config.has_value "mqtt_pass"; then
    password=$(bashio::config 'mqtt_pass')
  else
    bashio::log.info "No external MQTT service defined"
  fi
else
  host=$(bashio::services "mqtt" "host")
  password=$(bashio::services "mqtt" "password")
  port=$(bashio::services "mqtt" "port")
  username=$(bashio::services "mqtt" "username")
fi

ha_discovery_topic=$(bashio::config 'ha_discovery_topic')

echo "mqtt:" >/data/devices.yaml
echo "  broker: ${host}" >>/data/devices.yaml
echo "  port: ${port}" >>/data/devices.yaml
echo "  ha_discovery_topic: ${ha_discovery_topic}" >>/data/devices.yaml

# Process devices
echo "devices:" >>/data/devices.yaml

for device in $(bashio::config 'devices|keys'); do
  name=$(bashio::config "devices[${device}].name")
  echo "  - name: ${name}" >>/data/devices.yaml

  type=$(bashio::config "devices[${device}].type")
  echo "    type: ${type}" >>/data/devices.yaml

  model=$(bashio::config "devices[${device}].model")
  echo "    model: ${model}" >>/data/devices.yaml

  manufacturer=$(bashio::config "devices[${device}].manufacturer")
  echo "    manufacturer: ${manufacturer}" >>/data/devices.yaml

  if bashio::config.has_value devices[${device}].uid; then
    uid=$(bashio::config "devices[${device}].uid")
    echo "    uid: ${uid}" >>/data/devices.yaml
  fi

  if bashio::config.has_value devices[${device}].id; then
    id=$(bashio::config "devices[${device}].id")
    echo "    id: ${id}" >>/data/devices.yaml
  fi

  if bashio::config.has_value devices[${device}].channel; then
    channel=$(bashio::config "devices[${device}].channel")
    echo "    channel: ${channel}" >>/data/devices.yaml
  fi
done

/usr/bin/python3 /rattler/device_creator/main.py
