# snapserver/rootfs/etc/services.d/tts/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting TTS REST endpoint"
exec python3 /tts/main.py
