# rattler433/rootfs/etc/services.d/rtl_433/run
#!/usr/bin/with-contenv bashio
declare host
declare password
declare port
declare username
declare retain
declare output
declare prefix
declare customize_dir

if bashio::config.has_value "mqtt_host"; then
  if bashio::config.has_value "mqtt_host"; then
    host=$(bashio::config 'mqtt_host')
  else
    bashio::log.info "No external MQTT service defined"
  fi
  if bashio::config.has_value "mqtt_port"; then
    port=$(bashio::config 'mqtt_port')
  else
    bashio::log.info "No external MQTT service defined"
  fi
  if bashio::config.has_value "mqtt_user"; then
    username=$(bashio::config 'mqtt_user')
  else
    bashio::log.info "No external MQTT service defined"
  fi
  if bashio::config.has_value "mqtt_pass"; then
    password=$(bashio::config 'mqtt_pass')
  else
    bashio::log.info "No external MQTT service defined"
  fi
else
  host=$(bashio::services "mqtt" "host")
  password=$(bashio::services "mqtt" "password")
  port=$(bashio::services "mqtt" "port")
  username=$(bashio::services "mqtt" "username")
fi

if bashio::config.true 'retain'; then
  retain="1"
else
  retain="0"
fi

prefix=$(bashio::config 'mqtt_prefix')

output="output mqtt://${host}:${port},user=${username},pass=${password},retain=${retain},devices=${prefix}[/model][/channel][/id]"

echo ${output} >/rattler/output.conf

if bashio::config.true 'customize.active'; then
  customize_dir=$(bashio::config 'customize.folder')
  if ! bashio::fs.directory_exists '/share/rattler'; then
    mkdir -p /share/${customize_dir} && touch /share/${customize_dir}/flex.conf || bashio::exit.nok "Could not create Rattler Flex config folder!"
  fi
  if bashio::fs.file_exists "/share/${customize_dir}/flex.conf"; then
    bashio::log.log "Using flex.conf for custom flex devices..."
    rtl_433 -c /rattler/rtl_433.conf -c /rattler/sonoff_rm433.conf -c /share/${customize_dir}/flex.conf -c /rattler/output.conf
    else
    rtl_433 -c /rattler/rtl_433.conf -c /rattler/sonoff_rm433.conf -c /rattler/output.conf
  fi
  else
  rtl_433 -c /rattler/rtl_433.conf -c /rattler/sonoff_rm433.conf -c /rattler/output.conf
fi
