# snapserver/rootfs/etc/services.d/librespot1/run
#!/usr/bin/with-contenv bashio

if bashio::config.true 'Start Librespot?'; then
  bashio::log.info "Starting Librespot Server 1"

  exec /usr/bin/librespot \
  --name 'Snapcast Group 1' \
  --device-type avr \
  --backend pipe \
  --device /tmp/snap_librespot1 \
  --bitrate 160 \
  --disable-audio-cache \
  --enable-volume-normalisation \
  --volume-ctrl linear \
  --initial-volume=100 \
  >/dev/null 2>&1
fi
