# pwrstat/rootfs/etc/services.d/pwrstat/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs pwrstat
# ==============================================================================


bashio::log.info "Starting Pwrstat API..."

exec /pwrstat/pwrstat_api.py