# owntone/rootfs/etc/services.d/owntone/run
#!/usr/bin/with-contenv bashio

bashio::log.info 'Starting OwnTone Server...'
owntone -f -c /config/owntone/owntone.conf -P /var/run/owntone.pid
