# owntone/rootfs/etc/services.d/nginx/run
# #!/usr/bin/with-contenv bashio
# declare -a options

# # Wait for Forked-Daapd to become available
# bashio::net.wait_for 3689

# bashio::log.info "Starting Nginx..."

# exec nginx "${options[@]}"
