# snapserver/rootfs/etc/services.d/snapserver/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Snapcast server"
exec /usr/bin/snapserver
