# snapserver/rootfs/etc/services.d/snapserver/run
#!/usr/bin/with-contenv bashio
declare silence

silence=""
if bashio::config.true "silence_snapcast"; then
  silence=">/dev/null 2>&1"
  bashio::log.info "Silencing Snapcast"
fi

bashio::log.info "Starting Snapcast server"
exec /usr/bin/snapserver ${silence}
