# snapserver/rootfs/etc/services.d/snapserver/run
#!/usr/bin/with-contenv bashio
declare username
declare password
declare ip

if bashio::config.has_value "spotify_user"; then
  if bashio::config.has_value "spotify_pw"; then
    bashio::log.info "Setting up Spotify Config"

    username=$(bashio::config 'spotify_user')
    password=$(bashio::config 'spotify_pw')

    echo "stream = pipe:///config/snapfifo?name=TTS" >>/etc/snapserver.conf
    echo "stream = librespot:///librespot?name=Spotify &" >>/etc/snapserver.conf
    echo "dryout_ms=2000 &" >>/etc/snapserver.conf
    echo "username=${username} &" >>/etc/snapserver.conf
    echo "password=${password} &" >>/etc/snapserver.conf
    echo "devicename=Snapcast &" >>/etc/snapserver.conf
    echo "bitrate=320 &" >>/etc/snapserver.conf
    echo "wd_timeout=7800 &" >>/etc/snapserver.conf
    echo "volume=100 &" >>/etc/snapserver.conf
    echo "nomalize=true &" >>/etc/snapserver.conf
    echo "autoplay=false" >>/etc/snapserver.conf
    echo "stream = meta:///TTS/Spotify?name=Meta" >>/etc/snapserver.conf
  fi
fi
bashio::log.info "Starting Snapcast server"
snapserver
