# snapserver/rootfs/etc/services.d/mpd/run
#!/usr/bin/with-contenv bashio
/usr/bin/mpd --stdout --no-daemon /etc/mpd.conf
