# snapserver/rootfs/etc/services.d/mimic-py/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting TTS REST endpoint"
exec python3 /mimic-py/main.py
