# snapserver/rootfs/etc/services.d/librespot/run
#!/usr/bin/with-contenv bashio
declare name
declare pipe
declare silence
declare device_type
declare bitrate
declare additional_opts

if bashio::config.true "silence_librespot"; then
  silence=">/dev/null 2>&1"
fi

if bashio::config.true "start_librespot"; then
  for librespot in $(bashio::config 'librespot_instances|keys'); do

    name=$(bashio::config "librespot_instances[${librespot}].name")
    pipe=$(bashio::config "librespot_instances[${librespot}].pipe")

    if bashio::config.has_value 'device_type'; then
      device_type=$(bashio::config 'device_type')
    else
      device_type="speaker"
    fi

    if bashio::config.has_value 'bitrate'; then
      bitrate=$(bashio::config 'bitrate')
    else
      bitrate="320"
    fi

    if bashio::config.has_value 'additional_opts'; then
      additional_opts=$(bashio::config 'additional_opts')
    else
      additional_opts=""
    fi

    bashio::log.info "Starting ${name} Librespot at ${pipe}"
    exec /usr/bin/librespot \
    --name "${name}" \
    --device-type ${device_type} \
    --backend pipe \
    --device ${pipe} \
    --bitrate ${bitrate} \
    --disable-audio-cache \
    --enable-volume-normalisation \
    --volume-ctrl linear \
    --initial-volume=100 \
    ${additional_opts} \
    ${silence}
  done
fi
